# mkosi/systemd/wallet.service
[Unit]
Description=Cold Drive Wallet
After=weston.service

[Service]
Environment=XDG_RUNTIME_DIR=/run
ExecStart=/app/wallet
Restart=always

[Install]
WantedBy=graphical.target
